# deploy/dxf-generator.supervisor.conf
; ============================================================
; Supervisor config untuk DXF Generator (FastAPI / Uvicorn)
; ============================================================
; Install  : sudo cp dxf-generator.supervisor.conf /etc/supervisor/conf.d/
; Reload   : sudo supervisorctl reread && sudo supervisorctl update
; Status   : sudo supervisorctl status dxf-generator
; Restart  : sudo supervisorctl restart dxf-generator
; Log live : sudo tail -f /var/log/supervisor/dxf-generator.log
; ============================================================
;
; SESUAIKAN path di bawah sebelum digunakan:
;   /var/www/dxf-generator  →  path instalasi dxf-generator di server
;   www-data                →  user web server (www-data untuk Nginx/Apache)
; ============================================================

[program:dxf-generator]
command=/var/www/dxf-generator/.venv/bin/python run.py
directory=/var/www/dxf-generator
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/dxf-generator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/dxf-generator-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=www-data
; stopasgroup+killasgroup wajib agar semua worker process ikut berhenti
stopasgroup=true
killasgroup=true
; Beri waktu 30 detik untuk graceful shutdown (selesaikan request yang sedang jalan)
stopwaitsecs=30
