# deploy/laravel-queue.supervisor.conf
; ============================================================
; Supervisor config untuk Laravel Queue Worker
; ============================================================
; Install  : sudo cp laravel-queue.supervisor.conf /etc/supervisor/conf.d/
; Reload   : sudo supervisorctl reread && sudo supervisorctl update
; Status   : sudo supervisorctl status laravel-queue:*
; ============================================================
;
; SESUAIKAN path di bawah sebelum digunakan:
;   /var/www/aergas-apbn  →  path instalasi Laravel di server
;   www-data              →  user web server
; ============================================================

[program:laravel-queue]
command=php /var/www/aergas-apbn/artisan queue:work database --sleep=3 --tries=2 --max-time=3600 --memory=256
directory=/var/www/aergas-apbn
autostart=true
autorestart=true
startretries=3
; Jalankan 2 worker secara paralel untuk queue (EvidenceZip dll)
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/laravel-queue.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/laravel-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=www-data
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
